# packages/scheduler/eink.cron
WORKDIR=/home/marek/Projekty/vsb-eink-scheduler/examples
* * * * * all partial images/weather.html
5 * * * * all full images/lena.png
