# services/hoster/public/all.cron
WORKDIR=/home/marek/Projekty/vsb-eink-scheduler/examples
5-10 * * * * all full images/welcome.html
10-20 * * * * all partial images/weather.html
* * * * * all full images/lena.png images/baboon.jpg
